# curator4-1008298751.conf
[program:curator4-1008298751]
command=python3 /home/curator4/curator4-1008298751/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-1008298751/test.err.log
stdout_logfile=/home/curator4/curator4-1008298751/test.out.log